# automirror.timer
[Unit]
Description=Run automirror daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
